# systemd/fortigate-backup.service
[Unit]
Description=FortiGate Backup Portal
After=network.target

[Service]
Type=simple
User=fortigate-backup
WorkingDirectory=/opt/fortigate-backup
Environment=NODE_ENV=production
EnvironmentFile=/opt/fortigate-backup/.env
ExecStart=/usr/bin/pnpm start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
